# tests/nightly/cron/jet-user.cron
56 19 * * * /lfs3/projects/hfv3gfs/Samuel.Trahan/multi-app/tests/nightly/run-from-jet-user.sh tujet checkout > /lfs3/projects/hfv3gfs/Samuel.Trahan/scrub/tujet-nightly-checkout.log 2>&1

12 18 * * * rsync -arv /lfs3/projects/hfv3gfs/emc.nemspara/web/rt/. samuel.trahan@dmzgw.ncep.noaa.gov:/home/www/emc/htdocs/projects/rt/. > /lfs3/projects/hfv3gfs/Samuel.Trahan/scrub/nightly-rt-web-transfer.log 2>&1
