# tests/nightly/cron/theia-nemspara.cron
3 0 * * * mkdir -p /scratch4/NCEPDEV/nems/scrub/emc.nemspara

5 19 * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/nightly-test.log 2>&1
47 * * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/nightly-deliver.log 2>&1

5 19 * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia.intel test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/intel-nightly-test.log 2>&1
47 * * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia.intel deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/intel-nightly-deliver.log 2>&1

# Theia SLURM nightly test disabled for now because Theia SLURM thinks
# there are only 12 cores per node, thus doubling the number of nodes
# it takes to run the tests:
#
#  5 19 * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia.slurm.intel test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/slurm-nightly-test.log 2>&1
#  47 * * * * /scratch4/NCEPDEV/nems/noscrub/Samuel.Trahan/multi-app/tests/nightly/run-from-theia-nemspara.sh ngt theia.slurm.intel deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/slurm-nightly-deliver.log 2>&1
