# tests/nightly/cron/theia-nemspara.cron
30 01 * * * /home/emc.nemspara/bin/scrubber 1 3
25 13 * * * /home/emc.nemspara/bin/scrubber 2 3
30 07 * * * /home/emc.nemspara/bin/scrubber 3 3
25 19 * * * /home/emc.nemspara/bin/scrubber 4 3

# NEMS regression tests:
03 00 * * * mkdir -p /scratch4/NCEPDEV/nems/scrub/emc.nemspara

58 18 * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh ngt theia test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/nightly-test.log 2>&1
47  * * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh ngt theia deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/nightly-deliver.log 2>&1

58 18 * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh ngt theia.intel test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/intel-nightly-test.log 2>&1
47  * * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh ngt theia.intel deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/intel-nightly-deliver.log 2>&1

58 18 * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh -Mslurm ngt theia.slurm.intel test > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/slurm-nightly-test.log 2>&1
47  * * * * /scratch4/NCEPDEV/meso/noscrub/Ratko.Vasic/multi-app-test/tests/nightly/run-from-theia-nemspara.sh ngt theia.slurm.intel deliver > /scratch4/NCEPDEV/nems/scrub/emc.nemspara/slurm-nightly-deliver.log 2>&1
